# tools/procServMgr/procServMgr.service
[Unit]
Description=procServMgr - Check IOCs

[Service]
Type=simple
ExecStart=/usr/local/bin/procServMgr check

