# tools/procServMgr/procServMgr.service
[Unit]
Description=procServMgr - Check IOCs

[Service]
Type=simple
User=epics
Group=clas
ExecStart=/usr/local/bin/procServMgr check

