# tools/cron/clonsl3.crontab
# HEADER: This file was autogenerated at 2016-06-01 11:42:44 -0400 by puppet.
# HEADER: While it can still be managed manually, it is definitely not recommended.
# HEADER: Note particularly that the comments starting with 'Puppet Name' should
# HEADER: not be deleted, as doing so could cause duplicate cron jobs.
# Puppet Name: procServMgr_check_clas12
PATH=/bin:/usr/bin:/usr/clas12/release/pro/epics/bin
*/5 * * * * procServMgr check >> /usr/clas12/DATA/logs/`hostname -s`.log 2>&1
