# tools/cron/dc13.crontab
*/5 * * * * /usr/clas12/release/pro/epics/bin/procServMgr check
