# tools/cron/logs.crontab
0 4 * * * /usr/bin/find /usr/clas12/DATA/logs -mtime 30 -delete
