# tools/cron/procServ.crontab
PATH=/bin:/usr/bin:/usr/clas12/hps/prod/bin
*/5 * * * * procServMgr check >> /usr/clas12/hps/logs/`hostname -s`.log 2>&1
PATH=/bin:/usr/bin:/usr/clas12/release/pro/epics/bin
*/5 * * * * procServMgr check >> /usr/clas12/DATA/logs/`hostname -s`.log 2>&1
