# tools/cron/webopi.crontab
#
PATH=/bin:/usr/bin
#
*/5 * * * * rsync -azL --delete --exclude-from='/usr/clas12/release/pro/epics/css_share/webopi_excludes.txt' /usr/clas12/release/pro/epics/css_share/ /group/hallb/www/hallbopi/html/opis
#
