# tools/harp_fitter/harpdata2ops.cron
* * * * * /usr/clas12/release/pro/epics/tools/harp_fitter/harpdata2ops.sh >& /dev/null
